# vendor/cookbooks/poise-service/templates/default/systemd.service.erb
[Unit]
Description=<%= @name %>

[Service]
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
ExecReload=/bin/kill -<%= @reload_signal %> $MAINPID
KillSignal=<%= @stop_signal %>
User=<%= @user %>
WorkingDirectory=<%= @directory %>
Restart=<%= @restart_mode %>

[Install]
WantedBy=multi-user.target
