# vendor/cookbooks/poise-service/templates/default/sysvinit.sh.erb
#!/bin/sh
# Init script for <%= @name %> generated by poise-service
#
### BEGIN INIT INFO
# Provides: <%= @name %>
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Init script for <%= @name %>
# Description: Init script for <%= @name %>
### END INIT INFO

<%- @environment.each do |key, val| -%>
export <%= key %>="<%= val %>"
<%- end -%>

<%- if @platform_family == 'debian' -%>
. /lib/lsb/init-functions

_start() {
  start-stop-daemon --start --quiet --background \
      --pidfile "<%= @pid_file %>"<% unless @pid_file_external %> --make-pidfile<% end %> \
      --chuid "<%= @user %>" --chdir "<%= @directory %>" \
      --exec "<%= @daemon %>" -- <%= @daemon_options %>
}

_stop() {
  start-stop-daemon --stop --quiet --pidfile "<%= @pid_file %>" --user "<%= @user %>" --signal "<%= @stop_signal %>"
}

_status() {
  status_of_proc -p "<%= @pid_file %>" "<%= @daemon %>" "<%= @name %>"
}

_reload() {
  start-stop-daemon --stop --quiet --pidfile "<%= @pid_file %>" --user "<%= @user %>" --signal "<%= @reload_signal %>"
}

<%- elsif @platform_family == 'rhel' -%>
. /etc/rc.d/init.d/functions

_pid () {
  pidof -s -o $$ -o $PPID -o %PPID -x "$1" || pidof -s -o $$ -o $PPID -o %PPID -x "${1##*/}"
}

_start() {
  ( cd "<%= @directory %>" && daemon --user "<%= @user %>" --pidfile "<%= @pid_file %>" "<%= @daemon %> <%= @daemon_options %>" >/dev/null 2>&1 ) &
<%- unless @pid_file_external -%>
  sleep 1 # Give it some time to start before checking for a pid
  _pid "<%= @daemon %>" > "<%= @pid_file %>" || return 3
<%- end -%>
  return 0 # ¯\_(ツ)_/¯
}

_stop() {
  if [ -r "<%= @pid_file %>" ]; then
    kill -<%= @stop_signal%> "$(cat "<%= @pid_file %>")"
  else
    killproc "<%= @daemon %>" -<%= @stop_signal%> >/dev/null
  fi
}

_status() {
  status -p "<%= @pid_file %>" "<%= @daemon %>"
}

_reload() {
  if [ -r "<%= @pid_file %>" ]; then
    kill -<%= @reload_signal%> "$(cat "<%= @pid_file %>")"
  else
    return 1
  fi
}

<%- # Some functions to match LSB because RHEL doesn't (╯°□°）╯︵ ┻━┻ -%>

log_daemon_msg() {
  echo -n "$1"
}

log_progress_msg() {
  echo -n "$1"
}

log_warning_msg() {
  echo -n "$1"
}

log_failure_msg() {
  echo -n "$1"
}

log_end_msg() {
  if [ "$1" = 0 ]; then
    success
  else
    failure
  fi
}
<%- else -%>
<%- raise "Platform family #{@platform_family} is not supported" -%>
<%- end -%>

set -e

start() {
  if _start
  then
    rc=0
    sleep 1
    if ! kill -0 "$(cat "<%= @pid_file %>")" >/dev/null 2>&1; then
        log_failure_msg "<%= @name %> failed to start"
        rc=1
    fi
  else
    rc=1
  fi
  if [ "$rc" -eq 0 ]; then
    log_end_msg 0
  else
    log_end_msg 1
    rm -f "<%= @pid_file %>"
  fi
}

export PATH="${PATH:+$PATH:}/usr/sbin:/sbin"

case "$1" in
  start)
    log_daemon_msg "Starting <%= @name %>"
    if [ -s "<%= @pid_file %>" ] && kill -0 "$(cat "<%= @pid_file %>")" >/dev/null 2>&1; then
      log_progress_msg "apparently already running"
      log_end_msg 0
      exit 0
    fi
    start
  ;;

  stop)
    log_daemon_msg "Stopping <%= @name %>"
    _stop
    log_end_msg "$?"
    rm -f "<%= @pid_file %>"
  ;;

  reload|force-reload)
    log_daemon_msg "Reloading <%= @name %>"
    _reload
    log_end_msg "$?"
  ;;

  restart)
    set +e
    log_daemon_msg "Restarting <%= @name %>"
    if [ -s "<%= @pid_file %>" ] && kill -0 "$(cat "<%= @pid_file %>")" >/dev/null 2>&1; then
      _stop || true
      sleep 1
    else
      log_warning_msg "<%= @name %> not running, attempting to start."
      rm -f "<%= @pid_file %>"
    fi
    start
  ;;

  status)
    set +e
    _status
    exit $?
  ;;

  *)
    echo "Usage: /etc/init.d/<%= @name %> {start|stop|reload|force-reload|restart|status}"
    exit 1
esac

exit 0
